# motioneye/rootfs/etc/services.d/motioneye/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: motionEye
# Runs the motionEye
# ==============================================================================
declare -a options


bashio::log.info "Starting motionEye..."

options+=(-c /data/motioneye/motioneye.conf)

# Enable debug mode
if bashio::debug; then
    options+=(-d)
fi

# haaaaack
mount -t ext4 -o defaults,noauto,noatime,nodiratime /dev/sda1 /data1

# Run the motionEye
exec meyectl startserver "${options[@]}"
