# motioneye/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: motionEye
# Runs the Nginx daemon
# ==============================================================================

exec nginx -g "daemon off;"
